# systemd/user/load-mic.service
[Unit]
Description=Load external mic as ALSA source into PipeWire
After=pipewire.service pipewire-pulse.service wireplumber.service

[Service]
ExecStart=pactl load-module module-alsa-source device=hw:Generic_1,0
Type=oneshot

[Install]
WantedBy=default.target

